# systemd/.config/systemd/user/borgbackup.service
[Unit]
Description=Borg Backup Service

[Service]
EnvironmentFile=%h/.config/systemd/user/env.conf
Type=simple
ExecStart=/bin/bash %h/.borgbackup  

